# package/trafficmeter/files/traffic_meter.init
#!/bin/sh /etc/rc.common
START=85

start() {
	[ "X$($CONFIG get device_mode)" = "Xmodem" ] && return #do not start traffic meter on modem only.
	/sbin/cmd_traffic_meter start
}

stop() {
	/sbin/cmd_traffic_meter stop
}

boot()
{
	mkdir -p /tmp/traffic_meter
	start
}
