# package/telnetenable/files/telnet.init
#!/bin/sh /etc/rc.common

start()
{
	/usr/etc/functions/telnetenable-control start
}

stop()
{
	/usr/etc/functions/telnetenable-control stop
}

restart()
{
	stop
	start
}

boot()
{
	mkdir /dev/pts
	mknod -m 666 /dev/ptmx c 5 2
	mknod -m 666 /dev/pts/0 c 136 0
	mknod -m 666 /dev/pts/1 c 136 1

	start
}
